# com.user.downloads-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique LaunchAgent identifier -->
    <key>Label</key>
    <string>com.user.downloads-organizer</string>

    <!-- Command to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/.downloads-organizer/organize-downloads-ultra.sh --verbose 2>&1 | tee -a $HOME/.downloads-organizer/launchd.log</string>
    </array>

    <!-- Schedule for every Sunday at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <!-- 0 = Sunday, 1 = Monday, ..., 6 = Saturday -->
        <key>Weekday</key>
        <integer>0</integer>
        <!-- Hour (0-23) -->
        <key>Hour</key>
        <integer>10</integer>
        <!-- Minute (0-59) -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:$HOME/.npm-global/bin</string>
        <key>HOME</key>
        <string>$HOME</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>$HOME/.downloads-organizer</string>

    <!-- Output and error logs -->
    <key>StandardOutPath</key>
    <string>$HOME/.downloads-organizer/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/.downloads-organizer/stderr.log</string>

    <!-- Run on load (for initial test) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Keep trying if it fails -->
    <key>KeepAlive</key>
    <false/>

    <!-- Timeout in seconds (15 minutes) -->
    <key>TimeOut</key>
    <integer>900</integer>

    <!-- Low priority to not interfere -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Wait for network to be available (for Claude API) -->
    <key>WatchPaths</key>
    <array/>

</dict>
</plist>
